# com.chenxiaomei.ai-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chenxiaomei.ai-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/chenxiaomei/CodeBuddy/ai-daily/update_ai_daily.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>17</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/chenxiaomei/CodeBuddy/ai-daily</string>

  <key>StandardOutPath</key>
  <string>/Users/chenxiaomei/CodeBuddy/ai-daily/update.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chenxiaomei/CodeBuddy/ai-daily/update.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>GIT_TERMINAL_PROMPT</key>
    <string>0</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
